# debian/prerm
#!/bin/bash

if systemctl >/dev/null 2>&1; then
    systemctl --no-reload disable xlio.service >/dev/null 2>&1 || true
    systemctl stop xlio.service || true
fi
